# ops/txline-kit-dev-web.service.in
[Unit]
Description=TxLINE Kit development learning app
After=network-online.target

[Service]
Type=simple
User=claude
Group=claude
Environment=PORT=@PORT@
Environment=STATIC_ROOT=/srv/txline-kit/dev/current
ExecStart=/usr/bin/node /srv/txline-kit/static-server.mjs
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
